# deploy/kite_ws_price_fallback.service
[Unit]
Description=Auto Trader Kite WebSocket live-price fallback
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Auto_Trader
ExecStart=/bin/bash -lc 'cd /home/ubuntu/Auto_Trader && source /home/ubuntu/.autotrader_env && exec /home/ubuntu/Auto_Trader/venv/bin/python /home/ubuntu/Auto_Trader/scripts/kite_ws_price_fallback.py'
Restart=always
RestartSec=30
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
